# specs/javadocr.service
[Unit]
Description=javadocr

[Service]
ExecStart=/usr/bin/javadocr
Type=simple

[Install]
WantedBy=multi-user.target
